# .config/systemd/user/auto-backup.service
[Unit]
Description=Automatically do a backup

[Service]
Type=oneshot
ExecStart=/home/chasinglogic/.local/bin/auto-backup
