# .config/systemd/user/auto-backup.service
[Unit]
Description=Automatically do a backup

[Service]
ExecStart=%h/.local/bin/auto-backup %h
