# .config/systemd/user/mbsync.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=mbsync.service

[Install]
WantedBy=timers.target