# .config/systemd/user/regenerate_wiki_index.service
[Unit]
Description=VimWiki index regenerator service

[Service]
Type=oneshot
ExecStart=/home/chasinglogic/.local/bin/generate_vimwiki_index
